# packer/files/etc/systemd/system/docker.service.d/cluster.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --cluster-store=consul://ddd-01:8500 --cluster-advertise=eth1:2376
