# vagrant/nomad/nomad-client.service
[Unit]
Description=Nomad Client

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/nomad agent -config=/etc/nomad
