# vagrant/nomad/nomad-server.service
[Unit]
Description=Nomad Server

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/nomad agent -config=/etc/nomad
